# darkicecast/rootfs/etc/services.d/darkice/run
#!/usr/bin/with-contenv bashio
bashio::log.info "Starting DarkIce..."
exec /usr/bin/darkice -c /etc/darkice/darkice.cfg